# other/stream_server/build/com.stream.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">

    <!-- DO NOT EDIT UNLESS YOU KNOW WHAT YOU ARE DOING -->
    <dict>
        <key>Label</key>
        <!-- The label should be the same as the filename without the extension -->
        <string>com.stream.server</string>

        <key>WorkingDirectory</key>
        <string>/Users/williamjorgensen/Documents/01 Workplace William/stream_server/</string>

        <key>ProgramArguments</key>
        <array>
        <string>/usr/local/bin/node</string>
        <string>index.js</string>
        </array>

        <key>RunAtLoad</key>
        <true/>
        
        <key>KeepAlive</key>
        <true/>
    </dict>
</plist>